# contrib/seeder/opensyria-socat.service
[Unit]
Description=OpenSY DNS Socat Bridge (IPv4 to IPv6)
After=network.target opensy-seeder.service
Requires=opensy-seeder.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/socat UDP4-LISTEN:53,fork,reuseaddr UDP6:[::1]:5353
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
